# crontab
  */5 *  *   *   *     sh -c 'cd git_mirrors && ./update_mirrors.sh'
